# etc/seedbank/templates/seedbank_bootstrap
#!/bin/bash
### BEGIN INIT INFO
# Provides:          seedbank-bootstrap
# Required-Start:    $$all
# Required-Stop:
# Default-Start:     2
# Default-Stop:
# Short-Description: seedBank bootstrap script
# Description:       This init script is installed and generated by seedBank and will run disk
#                    scripts and Puppet recipes. It will erase itself during the first boot.
### END INIT INFO

# wait for network initialization
/bin/sleep 15

# run manifest(s) which are located in /seedbank
manifests () {
	/bin/mkdir /seedbank
	cd /seedbank
	/usr/bin/wget http://${server}:${port}/manifests.tgz -O /seedbank/manifests.tgz
	/bin/tar zxvf manifests.tgz
	export RUBYLIB=/seedbank/lib
	${manifests_commands}
}

# functions to run
${call_functions}

# cleanup
/usr/sbin/update-rc.d -f seedbank-bootstrap remove
/bin/rm /etc/init.d/seedbank-bootstrap
/bin/rm -rf /seedbank
	
exit 0
